# embitel/embitel_framework/conf/kewl.conf
[program:kewl-gunicorn]
user=root
environment=PATH='/home/embadmin/Documents/django_apps/django_experiments/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin'
command=/home/embadmin/Documents/django_apps/django_experiments/bin/gunicorn -c kewl_gunicorn.py  wsgi:application
directory=/home/embadmin/Documents/django_apps/django_experiments/embitel/embitel_framework/
autorestart=true
autostart=true
